# roles/postgresql10/tasks/files/postgresql-10.service
.include /lib/systemd/system/postgresql-10.service

[Service]
LimitNOFILE=81920
