# roles/pgbouncer/tasks/files/pgbouncer.service
.include /lib/systemd/system/pgbouncer.service

[Service]
LimitNOFILE=81920
